# embed/templates/systemd/webssh2.service
[Unit]
Description=foo

[Service]
WorkingDirectory=/opt/webssh2/app
ExecStart=npm start

[Install]
WantedBy=multi-user.target
